# etc/etc/systemd/system/docker.service.d/default_options.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -g /home/dany/docker
